# scripts/crontab
# run lotus-sync-restart every 15 min
*/15 * * * * /bin/lotus-sync-restart
